# root/etc/systemd/system/nut-server.service.d/override.conf
[Unit]
# Sometimes nut-server can fail to start at boot due to the network not being up
Wants=nut-driver.target network-online.target
Requires=network.target network-online.target
